# files/docker.service
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -c /etc/docker/daemon.json
